# pf-runtime/launchd/com.prettyfly.pf-triage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  PF Runtime communications triage — hourly cron.

  Fires at :05 past every hour to avoid top-of-hour congestion with other
  schedulers. Invokes `python -m pf_runtime comms triage --scheduled`,
  which iterates `with_mail_credentials()` accounts (calendar twins
  filtered), classifies new messages, writes proposals to the local
  SQLite store, and (Phase 2+) POSTs each proposal to the matching
  PFOS silo's agent-action endpoint.

  Install:
    cp pf-runtime/launchd/com.prettyfly.pf-triage.plist ~/Library/LaunchAgents/
    launchctl bootstrap "gui/$(id -u)" ~/Library/LaunchAgents/com.prettyfly.pf-triage.plist

  Unload:
    launchctl bootout "gui/$(id -u)" \
      ~/Library/LaunchAgents/com.prettyfly.pf-triage.plist

  Logs:
    ~/Library/Logs/pf-triage.out (JSON, one line per cycle)
    ~/Library/Logs/pf-triage.err
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.prettyfly.pf-triage</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/alexhale/Projects/agents/pf-runtime/.venv/bin/python</string>
    <string>-m</string>
    <string>pf_runtime</string>
    <string>comms</string>
    <string>triage</string>
    <string>--scheduled</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexhale/Projects/agents/pf-runtime</string>

  <!--
    StartCalendarInterval with only Minute=5 means "every hour at minute 5".
    Missing keys (Hour, Day, Month, Weekday) are treated as wildcards by launchd.
  -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Minute</key>
    <integer>5</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/alexhale/Library/Logs/pf-triage.out</string>

  <key>StandardErrorPath</key>
  <string>/Users/alexhale/Library/Logs/pf-triage.err</string>

  <!--
    Run only on the schedule, not at load time. We don't want the install
    command to immediately fire a cycle — operator runs `python -m pf_runtime
    comms triage` manually first to confirm registry + creds.
  -->
  <key>RunAtLoad</key>
  <false/>

  <!--
    Disable Spotlight-style abandonment: if the job runs long (90s budget),
    launchd shouldn't kill it before the cycle finishes.
  -->
  <key>ExitTimeOut</key>
  <integer>180</integer>
</dict>
</plist>
